# gke-credential-refresher.timer
[Unit]
Description=Refresh GKE credentials every 45 minutes
Requires=gke-credential-refresher.service

[Timer]
# Run every 45 minutes (tokens expire after 60 minutes)
OnBootSec=5min
OnUnitActiveSec=45min

[Install]
WantedBy=timers.target