# scripts/agents/bored/bored.service
[Unit]
Description=Bored commit queue agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=10
User=exedev
WorkingDirectory=/home/exedev/exe
ExecStart=/usr/bin/python3 /home/exedev/exe/scripts/agents/bored/bored.py
Environment=HOME=/home/exedev
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/home/exedev/.config/bored/env

[Install]
WantedBy=multi-user.target
