# exeuntu/welcome/welcome-exedev-webapp.service
[Unit]
Description=exe.dev box demo/welcome server
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/app
ExecStart=/app/welcome/welcomed
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
