# exeuntu/shelley.socket
[Unit]
Description=Shelley Socket

[Socket]
ListenStream=9999
Accept=no

[Install]
WantedBy=sockets.target
