# exeuntu/shelley.service
[Unit]
Description=Shelley AI Assistant Service
After=network.target
Wants=network.target

[Service]
Type=exec
User=exedev
Group=exedev
WorkingDirectory=/home/exedev
ExecStart=/usr/local/bin/shelley --config /exe.dev/shelley.json serve --port 9999
Restart=always
RestartSec=5
Environment=HOME=/home/exedev
Environment=USER=exedev
KillMode=process
TimeoutStopSec=10

# Generate gateway token before starting
ExecStartPre=+/bin/sh -c 'TOKEN=$(/usr/local/bin/generate-gateway-token) && echo "ANTHROPIC_API_KEY=$TOKEN" > /run/shelley-env && echo "OPENAI_API_KEY=$TOKEN" >> /run/shelley-env && echo "GEMINI_API_KEY=$TOKEN" >> /run/shelley-env'
EnvironmentFile=-/run/shelley-env

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=no
ReadWritePaths=/home/exedev

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shelley

[Install]
WantedBy=multi-user.target
