# exeuntu/shelley.service
[Unit]
Description=Shelley-the-Agent
After=network.target
Wants=network.target

[Service]
Type=exec
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/src
ExecStartPre=/bin/mkdir -p /home/exedev/.shelley
ExecStart=/usr/local/bin/shelley -debug -db /home/exedev/.shelley/shelley.db -config /exe.dev/shelley.json serve -port 9999
Restart=always
RestartSec=5
Environment=HOME=/home/exedev
Environment=USER=exedev
Environment=PATH=/headless-shell:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
KillMode=process
TimeoutStopSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=shelley

[Install]
WantedBy=multi-user.target
