# observability/woodpecker/woodpecker.timer
[Unit]
Description=Woodpecker daily observability timer

[Timer]
OnCalendar=*-*-* 07:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
