# ops/deploy/cgtop.service
[Unit]
Description=cgtop cgroup viewer
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'exec /usr/local/bin/cgtop -http "$(tailscale ip -4):9090"'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
