# ops/deploy/cgtop.service
[Unit]
Description=cgtop cgroup viewer
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/bin/bash -c 'exec /home/ubuntu/cgtop.latest -http "$(tailscale ip -4):9090"'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
