# ops/deploy/exelet-prod.service
[Unit]
Description=exe.dev exelet
After=network.target
Wants=network-online.target

[Service]
Type=simple
CPUWeight=1024
IOWeight=1024
LimitNOFILE=1048576
WorkingDirectory=/data/exelet
EnvironmentFile=-/etc/default/exelet
KillMode=process

# Use the latest exeletd
ExecStart=/bin/bash -c 'exec "/home/ubuntu/exeletd.latest" "-D" "--stage=prod" "--listen-address=tcp://$(tailscale ip -4):9080" "--http-addr=$(tailscale ip -4):9081"'

Restart=always
RestartSec=5

# Use journald for logging (remove file redirection)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
