# ops/deploy/exepipe-prod.service
[Unit]
Description=exe.dev exepipe
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu
EnvironmentFile=/etc/default/exepipe

# Use the latest timestamp version
ExecStart=/bin/bash -c 'exec "$(ls -t /home/ubuntu/exepipe.* | head -n1)" -controller -addr @exepipe -stage prod'

Restart=always
RestartSec=5

# Don't kill child processes. This is unusual but is the correct choice
# for exepipe. This lets us update exepipe without breaking existing
# user connections, which will continue to be handled by the old exepipe.
# The old exepipe will exit when all its connections are closed.
KillMode=process

# Use journald for logging (remove file redirection)
StandardOutput=journal
StandardError=journal

# Allow binding to privileged ports
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
