# ops/netactuate/roles/rc.local/templates/rc.local
#!/bin/sh -e
#
# rc.local
#
# Managed by ansible

if [ ! -f "/etc/ssh/ssh_host_dsa_key" ]
then
ssh-keygen -t dsa -N "" -f /etc/ssh/ssh_host_dsa_key
fi

if [ ! -f "/etc/ssh/ssh_host_rsa_key" ]
then
ssh-keygen -t rsa -N "" -f /etc/ssh/ssh_host_rsa_key
fi

{% if bgp_networks is defined %}
{%     if 'IPv4' in bgp_networks and bgp_networks['IPv4']|length %}
{%       for networkv4 in bgp_networks['IPv4'] if "/32" not in networkv4.net %}
ip -4 route add blackhole {{ networkv4.net }}&
{%       endfor %}
{%       for networkv4 in bgp_networks['IPv4'] %}
{%         if networkv4.bind_ranges is defined %}
{%           for range in networkv4.bind_ranges %}
{%             set parts = range.split('-') %}
{%             set base = parts[0].rsplit('.', 1) %}
for i in $(seq {{ base[1] }} {{ parts[1] }}); do ip -4 addr add local {{ base[0] }}.$i/32 dev lo 2>/dev/null & done
{%           endfor %}
{%         else %}
{%           for ipaddr in networkv4.ips %}
ip -4 addr add local {{ ipaddr }}/32 dev lo&
{%           endfor %}
{%         endif %}
{%       endfor %}
{%     endif %}
{%     if 'IPv6' in bgp_networks and bgp_networks['IPv6']|length %}
{%       for networkv6 in bgp_networks['IPv6'] if "/128" not in networkv6.net %}
ip -6 route add blackhole {{ networkv6.net }}&
{%       endfor %}
{%       for networkv6 in bgp_networks['IPv6'] %}
{%         for ipaddr in networkv6.ips %}
ip -6 addr add local {{ ipaddr }}/128 dev lo preferred_lft 0&
{%         endfor %}
{%       endfor %}
{%     endif %}
{% endif -%}

echo 1 > /proc/sys/net/ipv4/ip_forward
ping -c 1 8.8.8.8

exit 0
